# debian/urm.service
[Unit]
Description=URM Service

[Service]
Restart=on-failure
ExecStart=/usr/bin/urm
TimeoutStopSec=2

[Install]
WantedBy=multi-user.target
